# cookbooks/linux_web_stack/templates/default/adv_api.conf.erb
description "uWSGI server instance configured to serve adv_api"

start on runlevel [2345]
stop on runlevel [!2345]

respawn

setuid flask
setgid www-data

env PATH=/home/flask/venv/bin
chdir /home/flask
exec uwsgi --ini adv_api.ini